# role/node-exporter/templates/node-exporter.service
[Unit]
Description=Node Exporter

[Service]
ExecStart={{ node_exporter_bin }} --{{ exec_var }}

[Install]
WantedBy=multi-user.target